# proxy/ansible/roles/pipy/templates/pipy.service.j2
[Unit]
Description=pipy
Documentation=https://flomesh.io
After=network-online.target
Wants=network-online.target

[Service]
Environment=PIPY_CONFIG_FILE={% if role != "repo" %}
{{pipy_config_uri|default("/etc/pipy/proxy.js")}}
{% else %}
""
{% endif %}
Environment=PIPY_SPAWN={% if role != "repo" %}
{{ansible_processor_nproc - 1}}
{% else %}
0
{% endif %}
ExecStart=/usr/bin/pipy-start.sh pipy docker-start
ExecStop=/bin/killall pipy
User=pipy
WorkingDirectory=/etc/pipy
LimitNOFILE=655360

[Install]
WantedBy=multi-user.target
